# .github/workflows/bootiful-podcast-processor.service
[Unit]
Description=Bootiful Podcast Processor
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=ec2-user
ExecStart=/usr/bin/env pipenv run python /home/ec2-user/app/start.py

[Install]
WantedBy=multi-user.target
